# prometheus-grafana/py/crontab
*/15 * * * * python /get_data.py >> /photo/logfile
